# systemd/cursor-chat-export.timer
[Unit]
Description=Run cursor-chat-export periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
